# d10_sdk/D10_SDK_1028/ambarella/prebuild/third-party/armv7-a-hf/systemd/etc/systemd/system/mystartup.service
[Unit]
Description=11111111zwz startup script
[Service]
ExecStart=/home/web/zwzmount
[Install]
WantedBy=multi-user.target
